# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

env FOREMAN_PORT=<%= port %>
<% engine.env.each do |name,value| -%>
<% next if name.upcase == "FOREMAN_PORT" -%>
env <%= name %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

setuid <%= user %>

chdir <%= engine.root %>

exec <%= process.command %>
