# spec/resources/export/upstart/app-bravo-1.conf
start on starting app-bravo
stop on stopping app-bravo
respawn

env FOREMAN_PORT=5100

setuid app

chdir /tmp/app

exec ./bravo
